# fleaMarket/.cron
*/30 * * * * {cmd} "UserUploadTemp::cleanup()"