# webStudios/.cron
*/15 * * * * {cmd} checkAvailability
0 0 */1 * * {cmd} checkStudios
*/1 * * * * sudo /etc/init.d/nnweb-queue check
